# vm/provision/files/etc/supervisor/conf.d/mailcatcher.conf
[program:mailcatcher]
command=mailcatcher --ip 0.0.0.0 --no-quit --foreground
autorestart=true
user=nobody
group=nogroup
directory=/tmp
